# third_party/dbus-1.10.10/bus/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/home/patriotyk/sota_client_cpp/libdbus_build/var/run/dbus/system_bus_socket
